# wttrScript.timer
[Unit]
Description=Run wttrScript daily

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
